# goproject.conf
[supervisord]
logfile = /tmp/supervisord.log
[program:romanserver]
command=$GOPATH/bin/romanserver
autostart=true
autorestart=true
redirect_stderr=true